# infra/systemd/hhf-backup.timer
[Unit]
Description=Run hhf-backup nightly at 07:00 UTC

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
