# infra/systemd/hhf-backup.service
[Unit]
Description=Happy Hour Friends nightly pg_dump → S3
After=network-online.target postgresql-17.service

[Service]
Type=oneshot
User=hhf
ExecStart=/usr/local/bin/pg-backup.sh
